# aws/systemd/jiapp-stop-watchdog.service
[Unit]
Description=JiApp Stop Watchdog — checks idle and stops instance
After=jiapp.service

[Service]
Type=oneshot
ExecStart=/opt/jiapp/stop-watchdog.sh
StandardOutput=journal
StandardError=journal
